# pi-scripts/control-relay.service
[Unit]
Description=WebRTC DataChannel to UDP Control Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
SupplementaryGroups=dialout
WorkingDirectory=/home/pi

# Disable WiFi power saving before starting
# Using + prefix to run as root (required for iw command)
ExecStartPre=+/usr/sbin/iw wlan0 set power_save off

ExecStart=/usr/bin/python3 /home/pi/control-relay.py
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

# Car profile configuration (set to your tire/weight setup)
# Available profiles: badlands_4kg, ragnarok_stock
Environment=CAR_PROFILE=badlands_4kg

# Load secrets from environment file (create /home/pi/.env with TOKEN_SECRET)
# Can also override CAR_PROFILE in .env file
EnvironmentFile=-/home/pi/.env

# CPU priority for low-latency control
Nice=-10
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=50

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=control-relay

[Install]
WantedBy=multi-user.target
